# Bots/ytbot/etc/systemd/system/telegram-bot-api.service
[Unit]
Description=Telegram Bot API (Local)
After=network.target

[Service]
User=typezero
WorkingDirectory=/mnt/nvme1/work/telegram-bot-api
ExecStart=/home/typezero/src/telegram-bot-api/build/telegram-bot-api \
  --api-id YOUR_API_ID \
  --api-hash YOUR_API_HASH \
  --local \
  --http-port 8081 \
  --dir /mnt/nvme1/work/telegram-bot-api

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
